# context-synthesizer/packaging/DEBIAN/postinst
#!/bin/bash
set -e

CS_USER="context-synthesizer"
CS_HOME="/var/lib/context-synthesizer"
CS_LIB="/usr/lib/context-synthesizer"
CS_STATS="${CS_HOME}/stats"

case "${1}" in
  configure)
    if ! id "$CS_USER" &>/dev/null; then
      adduser --system --group --home "$CS_HOME" --no-create-home "$CS_USER"
    fi

    mkdir -p "$CS_STATS"
    chown -R "${CS_USER}:${CS_USER}" "$CS_HOME"
    chmod 750 "$CS_HOME" "$CS_STATS"

    if [[ ! -d "${CS_LIB}/.venv" ]]; then
      python3 -m venv "${CS_LIB}/.venv"
    fi
    "${CS_LIB}/.venv/bin/pip" install --quiet --upgrade pip
    "${CS_LIB}/.venv/bin/pip" install --quiet -r "${CS_LIB}/requirements.txt"

    # Default developer id to the user who installed the package
    if [[ -f /etc/default/context-synthesizer ]]; then
      if ! grep -q '^TELEMETRY_DEVELOPER_ID=.' /etc/default/context-synthesizer; then
        INSTALLER="${SUDO_USER:-${USER:-unknown}}"
        sed -i "s/^TELEMETRY_DEVELOPER_ID=.*/TELEMETRY_DEVELOPER_ID=${INSTALLER}/" \
          /etc/default/context-synthesizer || true
      fi
    fi

    systemctl daemon-reload || true
    if systemctl is-enabled context-synthesizer &>/dev/null; then
      systemctl restart context-synthesizer || true
    fi

    echo ""
    echo "Context Synthesizer installed."
    echo ""
    echo "Claude CLI developers (BYOK — recommended):"
    echo "  1. context-synthesizer-setup-user"
    echo "  2. ~/.claude/settings.json → ANTHROPIC_BASE_URL + ANTHROPIC_API_KEY"
    echo "  3. Stats: ~/.local/share/context-synthesizer/stats/telemetry.jsonl"
    echo ""
    echo "Shared system service (optional, team key in /etc/default):"
    echo "  sudo systemctl enable --now context-synthesizer"
    echo "  Telemetry: ${CS_STATS}/telemetry.jsonl"
    echo ""
    echo "  Docs: /usr/share/doc/context-synthesizer/Usage.md"
    echo ""
    ;;
esac

exit 0
