# context-synthesizer/packaging/DEBIAN/prerm
#!/bin/bash
set -e

case "${1}" in
  remove|deconfigure)
    if systemctl is-active context-synthesizer &>/dev/null; then
      systemctl stop context-synthesizer || true
    fi
    if systemctl is-enabled context-synthesizer &>/dev/null; then
      systemctl disable context-synthesizer || true
    fi
    ;;
esac

exit 0
